# data/greeter-autostart/gnome-settings-daemon.desktop.in
[Desktop Entry]
Type=Application
_Name=GNOME Settings Daemon
TryExec=/usr/libexec/gnome-settings-daemon
Exec=/usr/libexec/gnome-settings-daemon --gconf-prefix=/apps/gdm/simple-greeter/settings-manager-plugins
OnlyShowIn=GNOME;
X-GNOME-Autostart-Phase=Initialization
X-GNOME-Autostart-Notify=true
X-GNOME-AutoRestart=true
